# .docker/chrome/chromium-cdp.conf
[program:chromium-cdp]
priority=2
user=seluser
environment=DISPLAY=":99.0",HOME="/home/seluser"

; Start a persistent Chromium instance for noVNC + CDP usage
; Exposes DevTools endpoint at http://<host>:9222/json/version
command=/bin/bash -lc "set -e; PROFILE_DIR=/home/seluser/chrome-profile; mkdir -p \"$PROFILE_DIR\"; rm -f \"$PROFILE_DIR\"/SingletonLock \"$PROFILE_DIR\"/SingletonCookie \"$PROFILE_DIR\"/SingletonSocket || true; exec /usr/bin/chromium --disable-dev-shm-usage --disable-gpu --remote-debugging-address=0.0.0.0 --remote-debugging-port=9223 --remote-allow-origins=* --user-data-dir=\"$PROFILE_DIR\" --window-size=1138,712 about:blank"

autostart=true
autorestart=true
startsecs=1
startretries=10
stopasgroup=true
killasgroup=true

redirect_stderr=true
stdout_logfile=/var/log/supervisor/chromium-cdp.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5


[program:cdp-forwarder]
priority=3
user=seluser
environment=LISTEN_HOST="0.0.0.0",LISTEN_PORT="9222",TARGET_HOST="127.0.0.1",TARGET_PORT="9223"
command=/usr/bin/python3 /opt/cdp/tcp_forwarder.py
autostart=true
autorestart=true
startsecs=1
startretries=10
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cdp-forwarder.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3

